# .config/systemd/user/cron-dd543bb81b62.timer
[Unit]
Description=cron-dd543bb81b62 timer unit

[Timer]
OnCalendar=*:*
