# .config/systemd/user/cron-dd543bb81b62.service
[Unit]
Description=cron-dd543bb81b62 service unit

[Service]
TimeoutStartSec=0
ExecStart=curl $(cat ~/.uptime-url)
Type=oneshot
