# .config/systemd/user/cron-53dc9761798a.timer
[Unit]
Description=cron-53dc9761798a timer unit

[Timer]
OnCalendar=*:*
