# .config/systemd/user/cron-53dc9761798a.service
[Unit]
Description=cron-53dc9761798a service unit

[Service]
TimeoutStartSec=0
ExecStart=bash update_git.sh
Type=oneshot
