# .config/systemd/user/cron-d4e24025f0a5.timer
[Unit]
Description=cron-d4e24025f0a5 timer unit

[Timer]
OnCalendar=*:*
