# .config/systemd/user/cron-dd543bb81b62.service
[Unit]
Description=Uptime 

[Service]
TimeoutStartSec=0
ExecStart=curl $(cat %h/.uptime-url)
Type=oneshot
