# .config/systemd/user/uptime.service
[Unit]
Description=Uptime 

[Service]
TimeoutStartSec=0
ExecStart=curl $(cat /home/neon/.uptime-url)
Type=oneshot
