# .config/systemd/user/uptime.timer
[Unit]
Description=Uptime for 1min

[Timer]
OnCalendar=*:*
